# SteelHead/KVM-Cookbooks/101-virtual-fail-to-wire/etc/systemd/system/cdp-watchdog.service
[Unit]
Description=CDP Watchdog
[Service]
Type=simple
ExecStart=/usr/sbin/cdp-watchdog.sh
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target
